# iso-builder/files/blackgate.service
[Unit]
Description=Blackgate SRT Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=blackgate
Group=blackgate
WorkingDirectory=/opt/blackgate

Environment=PHX_SERVER=true
Environment=PORT=4000
Environment=PHX_HOST=0.0.0.0
Environment=API_AUTH_USERNAME=admin
Environment=API_AUTH_PASSWORD=password123
Environment=LICENSE_SERVER_URL=https://license-server-eta-bay.vercel.app
Environment=HOME=/opt/blackgate

ExecStart=/opt/blackgate/bin/blackgate start
ExecStop=/opt/blackgate/bin/blackgate stop

Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
