# iso-builder/config/includes.chroot/etc/systemd/system/docker.service.d/override.conf
[Unit]
After=var-lib-docker.mount network.target containerd.service
Requires=var-lib-docker.mount
